# systemd/multi-tailnet-setup.service
[Unit]
Description=Multi-tailnet setup (netns, routes, firewall)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/multi-tailnet.sh setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
